# com.unit-tracker.hai-evaluator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.unit-tracker.hai-evaluator</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/wycliff/Documents/unit_tracker/.venv/bin/python</string>
        <string>/Users/wycliff/Documents/unit_tracker/main.py</string>
        <string>hai_evaluator</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/wycliff/Documents/unit_tracker</string>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>/Users/wycliff/Documents/unit_tracker/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wycliff/Documents/unit_tracker/launchd.log</string>
</dict>
</plist>
